# backups/2026-04-22_134803/meta/launchd/plists/ai.openclaw.idea-proactive-curator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.idea-proactive-curator</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>$HOME/.openclaw/idea/lead/workflow/curate-proactive-knowledge.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/thomas</string>
    <key>PATH</key>
    <string>/Users/thomas/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/thomas</string>
  <key>StandardOutPath</key>
  <string>/Users/thomas/.openclaw/idea/lead/workflow/logs/proactive-curator.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thomas/.openclaw/idea/lead/workflow/logs/proactive-curator.err.log</string>
</dict>
</plist>
