# backups/2026-04-22_134803/meta/launchd/plists/ai.openclaw.command-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.command-listener</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>$HOME/.openclaw-rescue/command-listener.sh &gt;&gt; $HOME/.openclaw-rescue/command-listener.log 2&gt;&amp;1</string>
  </array>
  <key>StartInterval</key>
  <integer>20</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/thomas</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/thomas/.openclaw-rescue/command-listener.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thomas/.openclaw-rescue/command-listener.log</string>
</dict>
</plist>
