# backups/2026-04-22_134803/meta/launchd/plists/com.openclaw.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>PATH="/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin" $HOME/openclaw-backup/backup.sh &gt;&gt; $HOME/.openclaw-backup-keys/backup.log 2&gt;&amp;1</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/openclaw-backup.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/openclaw-backup.log</string>
  <key>StartInterval</key>
  <integer>21600</integer>
</dict>
</plist>
