# backups/2026-04-22_195056/meta/launchd/plists/ai.openclaw.idea-lead-project-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.idea-lead-project-watch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>$HOME/.openclaw/idea/lead/workflow/project-watch.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>20</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/thomas</string>
    <key>PATH</key>
    <string>/Users/thomas/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/thomas/.openclaw/idea/lead/workflow/logs/project-watch.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thomas/.openclaw/idea/lead/workflow/logs/project-watch.launchd.err.log</string>
</dict>
</plist>