# backups/2026-04-22_074414/meta/launchd/plists/ai.openclaw.pm01-daily-writer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.pm01-daily-writer</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>$HOME/.openclaw/workspace-pm01/workflow/auto-daily.sh &gt;&gt; $HOME/.openclaw/workspace-pm01/workflow/logs/auto-daily.launchd.log 2&gt;&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/thomas</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/thomas/.openclaw/workspace-pm01/workflow/logs/auto-daily.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/thomas/.openclaw/workspace-pm01/workflow/logs/auto-daily.launchd.log</string>
</dict>
</plist>
