# backups/2026-04-21_003522/meta/launchd/plists/ai.openclaw.oa-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.oa-dashboard</string>

  <key>ProgramArguments</key>
	<array>
		<string>/Users/thomas/.openclaw/tools/oa-cli/.venv/bin/oa</string>
		<string>serve</string>
		<string>--host</string>
		<string>0.0.0.0</string>
		<string>--port</string>
		<string>3460</string>
		<string>--no-open</string>
	</array>

  <key>WorkingDirectory</key>
  <string>/Users/thomas/.openclaw/tools/oa-cli/oa-project</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/thomas/.openclaw/logs/oa-dashboard.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/thomas/.openclaw/logs/oa-dashboard.err.log</string>
</dict>
</plist>
